# scripts/deploy/com.agentlab.trading-watcher.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.agentlab.trading-watcher</string>
  <key>ProgramArguments</key>
  <array>
    <string>__AGENT_LAB_VENV_PYTHON__</string>
    <string>__AGENT_LAB_ROOT__/scripts/run_trading_mission_watcher.py</string>
  </array>
  <key>EnvironmentVariables</key>
  <dict>
    <key>QUANT_PIPELINE_ROOT</key>
    <string>__QUANT_PIPELINE_ROOT__</string>
    <key>AGENTIC_TRADING_DB</key>
    <string>__AGENTIC_TRADING_DB__</string>
    <key>AGENT_LAB_TRADING_WATCHER_COOLDOWN_SEC</key>
    <string>600</string>
  </dict>
  <key>StartInterval</key>
  <integer>600</integer>
  <key>StandardOutPath</key>
  <string>__AGENT_LAB_HOME__/logs/trading-watcher.log</string>
  <key>StandardErrorPath</key>
  <string>__AGENT_LAB_HOME__/logs/trading-watcher.err</string>
</dict>
</plist>
